# deploy/fingerprint-sns.service
[Unit]
Description=Fingerprint System SNS - Flask Application
After=network.target mysql.service

[Service]
Type=simple
User=fingerprint
Group=fingerprint
WorkingDirectory=/opt/fingerprint-sns
Environment="PATH=/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONUNBUFFERED=1"
ExecStart=/usr/bin/python3 /opt/fingerprint-sns/run_production.py
Restart=always
RestartSec=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
